# vllm.service
[Unit]
Description=vLLM Model Server (GLM-OCR)
After=network.target

[Service]
User=sysadmin
Group=sysadmin
WorkingDirectory=/home/sysadmin/GLM-OCR
Environment="PATH=/home/sysadmin/GLM-OCR/.venv/bin:/usr/bin:bin"
Environment="CC=/usr/bin/gcc"
Environment="CXX=/usr/bin/g++"
Environment="http_proxy=http://172.29.111.242:7890"
Environment="https_proxy=http://172.29.111.242:7890"
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="VLLM_USE_V1=0"
ExecStart=/home/sysadmin/GLM-OCR/.venv/bin/vllm serve zai-org/GLM-OCR --port 8001 --trust-remote-code --tensor-parallel-size 1 --gpu-memory-utilization 0.85 --speculative-config.method mtp --speculative-config.num_speculative_tokens 1 --enable-prefix-caching
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target